# aliminiapp.conf
[program:aliminiapp]
command = gunicorn aliminiapp.wsgi_prod -b 127.0.0.1:8032
directory = /root/AlipayMiniApp/aliminiapp/
autostart = true
autorestart = true